# conf/openvpn@.service
.include /lib/systemd/system/openvpn@.service

[Service]
Restart=always
RestartSec=1
